# service/wordle_refresh.timer.template
[Unit]
Description=Refresh the wordle solution
RefuseManualStart=no
RefuseManualStop=no

[Timer]
Persistent=true
OnCalendar=daily
Unit=wordle.service 

[Install]
WantedBy=timers.target
